# supervisord.conf
[supervisord]
nodaemon=true

[program:t38modem]
command=/bin/bash -c '/usr/bin/t38modem -p +/dev/ttyT0 --force-fax-mode --no-force-t38-mode --no-h323 --sip-register $SIPUSER@$SIPADDR,$SIPPASS --sip-listen udp$*:5060 --route "modem:.*=sip:<dn>@$SIPADDR" --route "sip:.*=modem:<dn>"'
autostart=true
autorestart=true
post_start_hook=/bin/bash -c "chown uucp:uucp /dev/ttyT0"

[program:hylafax]
command=/bin/bash '/usr/sbin/hylafax-init.sh'
autostart=true
